# contrib/initscripts/mygpiod.service.in
#
# SPDX-License-Identifier: GPL-3.0-or-later
# myGPIOd (c) 2020-2023 Juergen Mang <mail@jcgames.de>
# https://github.com/jcorporation/mygpiod
#
[Unit]
Description=myGPIOd server daemon
Requires=local-fs.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/mygpiod
Type=simple

User=mygpiod
Group=gpio

RuntimeDirectory=mygpiod

# disallow writing to /usr, /bin, /sbin, ...
ProtectSystem=yes

# more paranoid security settings
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes

[Install]
WantedBy=multi-user.target
