# contrib/initscripts/mygpiod.openrc.in
#!/sbin/openrc-run
#
# SPDX-License-Identifier: GPL-3.0-or-later
# myGPIOd (c) 2020-2022 Juergen Mang <mail@jcgames.de>
# https://github.com/jcorporation/myGPIOd

name=myGPIOd
command="@CMAKE_INSTALL_FULL_BINDIR@/mygpiod"
command_args=""
command_background="yes"
command_user="mygpiod:gpio"

start_stop_daemon_args=""
pidfile="/run/$name.pid"

depend() {
	after firewall
}
