# contrib/initscripts/mygpiod.service.in
# SPDX-License-Identifier: GPL-3.0-or-later
# myMPD (c) 2020-2025 Juergen Mang <mail@jcgames.de>
# https://github.com/jcorporation/mygpiod

[Unit]
Description=myGPIOd server daemon
Documentation=man:mygpiod(1)
Requires=local-fs.target

[Service]
DynamicUser=yes
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/mygpiod
LockPersonality=yes
PrivateDevices=no
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectProc=invisible
RestrictRealtime=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX
RestrictNamespaces=yes
RuntimeDirectory=mygpiod
RuntimeDirectoryMode=0750
SystemCallArchitectures=native
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
